# LaunchAgent/com.jamesonmorrill.vigilcore.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jamesonmorrill.vigilcore</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/jamesonmorrill/Desktop/Local AI/dailybrief/vigil-core/dist/index.js</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/jamesonmorrill/Desktop/Local AI/dailybrief/vigil-core</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>
    <key>ProcessType</key>
    <string>Background</string>
    <key>StandardOutPath</key>
    <string>/Users/jamesonmorrill/Library/Logs/DailyBrief/vigilcore-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jamesonmorrill/Library/Logs/DailyBrief/vigilcore-stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>ANTHROPIC_API_KEY</key>
        <string>REDACTED_API_KEY</string>
        <key>PORT</key>
        <string>3001</string>
    </dict>
</dict>
</plist>
